# mordzix.service
[Unit]
Description=Mordzix PRO ULTRA Server
After=network.target
Wants=network-online.target
StartLimitBurst=0

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/mordzix
Environment=PATH=/opt/mordzix/.venv/bin
ExecStart=/opt/mordzix/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 6969 --loop asyncio
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutSec=30

# Restart configuration
Restart=always
RestartSec=5

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/mordzix/data /opt/mordzix/logs
DynamicUser=no

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mordzix

[Install]
WantedBy=multi-user.target